# debian/ugreen-dxp-leds-dkms.service
[Unit]
Description=Bind UGREEN DXP LED I2C controller
Documentation=https://github.com/lpaolini/ugreen-dxp-leds-dkms
After=systemd-modules-load.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/ugreen-dxp-leds-dkms/probe-led-ugreen
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
